# systemd/habit-bot.service
[Unit]
Description=Discord Habit Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=habitbot
Group=habitbot
WorkingDirectory=/opt/discord-habit-bot
Environment=PATH=/opt/miniconda3/envs/habit-bot/bin
ExecStart=/opt/miniconda3/envs/habit-bot/bin/python /opt/discord-habit-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/discord-habit-bot

# Environment file
EnvironmentFile=/opt/discord-habit-bot/.env

[Install]
WantedBy=multi-user.target