# systemd/habit-bot.service
[Unit]
Description=Discord Habit Bot - Gamified Habit Tracking
After=network.target
Wants=network.target

[Service]
Type=simple
User=ian
Group=ian
WorkingDirectory=/home/ian/scripts/discord-buddy
Environment=PATH=/home/ian/miniconda3/envs/habit-bot/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -c 'source /home/ian/miniconda3/etc/profile.d/conda.sh && conda activate habit-bot'
ExecStart=/home/ian/miniconda3/envs/habit-bot/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ian/scripts/discord-buddy
ProtectHome=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target